# supervisord.conf
[program:aroundtheclock]
directory=<PATH_TO_PROJECT_ROOT>
command=<PATH_TO_PROJECT_ROOT>/run.sh
autostart=true
autorestart=true
